# docs/raspberry_pi_config/conf/99-yubikey-remove.rules
ACTION=="remove", SUBSYSTEM=="usb", ENV{PRODUCT}=="1050/407/571", RUN+="/usr/local/bin/start-mission-script.sh"
